# nvidia-template/debian/nvidia-assemble.service.in
[Unit]
Description=Service for linking and adding signature to nvidia driver in Live system.
Before=basic.target systemd-modules-load.service systemd-udevd.service
Before=nvidia-persistenced.service lernstick-xserver-xorg.service
Conflicts=shutdown.target
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/nvidia-assemble-${MODULE_VERSION} \
    "/usr/lib/nvidia/${MODULE_NAME}-${MODULE_VERSION}/${KERNEL_VERSION}" \
    "/lib/modules/${KERNEL_VERSION}/extra/${MODULE_NAME}-${MODULE_VERSION}" \
    "/usr/lib/linux-kbuild-${KBUILD_VERSION}/scripts/sign-file" \
    "${KERNEL_VERSION}"
ExecStop=/usr/bin/rm -rf /lib/modules/${KERNEL_VERSION}/extra/${MODULE_NAME}-${MODULE_VERSION}

[Install]
RequiredBy=basic.target
